# sync/tablet/systemd/remarkable-push-sync.timer
[Unit]
Description=reMarkable push sync backstop (real triggers: edit/sleep/wake)

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
AccuracySec=5min
Unit=remarkable-push-sync.service

[Install]
WantedBy=timers.target
