# alt-ui/sync/server/systemd/alt-ui-upload.service
[Unit]
Description=alt-ui inbound upload + PDF render service
After=network.target

[Service]
ExecStart=/usr/bin/node /home/exedev/bin/alt-ui-upload.js
Restart=always
RestartSec=2
User=exedev
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
